# goss.timer
[Unit]
Description=Goss test

[Timer]
OnUnitActiveSec=5m
OnBootSec=1m

[Install]
WantedBy=timers.target
